# Supervisorfile.ini
[program:node]
command=/nodejs/bin/node /var/www/openurban/server.js
numprocs=1
directory=/var/www/openurban/
redirect_stderr=true
stdout_logfile=/var/log/openurban/openurban.log
autostart=true
autorestart=true
minfds=100000
strip_ansi=true

[inet_http_server]
port=9001


[supervisord]
logfile = /var/log/supervisord.log
